# CryptoIntelNode_workspace/deploy/systemd/crypto-intel-api.service
[Unit]
Description=Crypto Intel Node API
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=crypto-intel
Group=crypto-intel
WorkingDirectory=/opt/crypto-intel-node/current
EnvironmentFile=/home/crypto-intel/.config/crypto-intel-node/runtime.env
ExecStart=/usr/bin/node /opt/crypto-intel-node/current/scripts/start-api.js
Restart=on-failure
RestartSec=5
UMask=0077
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/crypto-intel/.local/state/crypto-intel-node

[Install]
WantedBy=multi-user.target
