# CryptoIntelNode_workspace/deploy/systemd/crypto-intel-a2a-worker.path
[Unit]
Description=Watch Crypto Intel Node A2A ready jobs

[Path]
PathExistsGlob=/home/crypto-intel/.local/state/crypto-intel-node/a2a/jobs/*/ready
Unit=crypto-intel-a2a-worker.service

[Install]
WantedBy=multi-user.target
